# psp/filesystem/usr/lib/opkg/info/libgles-omap3.list
/usr/lib/libEGL.so
/usr/lib/libGLES_CM.so.1.6.16.3977
/usr/lib/libOpenVGU.so.1.6.16.3977
/usr/lib/ES3.0/libOpenVG.so.1.6.16.3977
/usr/lib/libIMGegl.so.1.6
/usr/bin/ES3.0/pvrsrvinit
/usr/lib/ES5.0/libGLESv2.so.1.6.16.3977
/usr/lib/libPVRScopeServices.so
/usr/bin/ES2.0/pvr2d_test
/usr/lib/libsrv_init.so
/usr/lib/libsrv_init.so.1.6
/usr/lib/ES2.0/libGLES_CM.so.1.6.16.3977
/usr/lib/ES2.0/libOpenVGU.so.1.6.16.3977
/usr/lib/libOpenVGU.so
/usr/lib/libsrv_init.so.1
/usr/lib/libGLESv2.so.1.6.16.3977
/usr/lib/libIMGegl.so
/usr/lib/ES2.0/libOpenVG.so.1.6.16.3977
/usr/lib/ES6.0/libsrv_init.so.1.6.16.3977
/usr/lib/libpvr2d.so
/usr/lib/ES6.0/libIMGegl.so.1.6.16.3977
/usr/lib/libpvr2d.so.1.6
/usr/lib/libpvr2d.so.1
/usr/lib/libpvrPVR2D_BLITWSEGL.so
/usr/lib/libGLESv2.so.1.6
/usr/lib/ES6.0/libsrv_um.so.1.6.16.3977
/usr/lib/libEGL.so.1.6
/usr/bin/ES6.0/pvrsrvinit
/etc/init.d/pvr-init
/usr/lib/ES5.0/libsrv_init.so.1.6.16.3977
/usr/lib/libusc.so
/usr/bin/ES5.0/pvr2d_test
/usr/lib/ES3.0/libpvr2d.so.1.6.16.3977
/usr/lib/ES5.0/libIMGegl.so.1.6.16.3977
/etc/powervr.ini
/usr/lib/libGLES_CM.so.1.6
/usr/lib/libsrv_um.so.1.6
/usr/lib/ES3.0/libGLESv2.so.1.6.16.3977
/usr/lib/ES3.0/libEGL.so.1.6.16.3977
/usr/lib/libpvrPVR2D_FRONTWSEGL.so
/usr/lib/ES5.0/libsrv_um.so.1.6.16.3977
/usr/lib/libEGL.so.1.6.16.3977
/usr/lib/libglslcompiler.so.1.6
/usr/lib/libIMGegl.so.1.6.16.3977
/usr/lib/ES3.0/libsrv_init.so.1.6.16.3977
/usr/bin/ES5.0/pdump
/usr/lib/libpvrPVR2D_LINUXFBWSEGL.so
/usr/bin/ES3.0/pvr2d_test
/usr/lib/ES2.0/libpvr2d.so.1.6.16.3977
/usr/lib/libsrv_um.so.1
/usr/lib/libsrv_um.so.1.6.16.3977
/usr/lib/libOpenVGU.so.1
/usr/lib/ES2.0/libglslcompiler.so.1.6.16.3977
/usr/lib/ES2.0/libGLESv2.so.1.6.16.3977
/usr/lib/libglslcompiler.so.1
/usr/lib/ES3.0/libglslcompiler.so.1.6.16.3977
/usr/lib/ES2.0/libsrv_init.so.1.6.16.3977
/usr/lib/ES6.0/libOpenVG.so.1.6.16.3977
/usr/lib/libglslcompiler.so
/usr/lib/ES6.0/libGLES_CM.so.1.6.16.3977
/usr/lib/libsrv_um.so
/usr/lib/ES6.0/libOpenVGU.so.1.6.16.3977
/usr/lib/ES5.0/libglslcompiler.so.1.6.16.3977
/usr/lib/ES2.0/libEGL.so.1.6.16.3977
/usr/lib/libpvrPVR2D_FLIPWSEGL.so
/usr/lib/libGLESv2.so
/usr/lib/ES6.0/libEGL.so.1.6.16.3977
/usr/lib/libIMGegl.so.1
/usr/lib/libGLES_CM.so.1
/usr/bin/ES2.0/pvrsrvinit
/usr/lib/ES6.0/libglslcompiler.so.1.6.16.3977
/usr/lib/ES3.0/libIMGegl.so.1.6.16.3977
/usr/lib/libEGL.so.1
/usr/bin/ES6.0/pvr2d_test
/usr/bin/cputype
/usr/lib/ES5.0/libOpenVG.so.1.6.16.3977
/usr/lib/ES3.0/libsrv_um.so.1.6.16.3977
/usr/lib/ES5.0/libGLES_CM.so.1.6.16.3977
/usr/bin/ES6.0/pdump
/usr/lib/ES5.0/libOpenVGU.so.1.6.16.3977
/usr/lib/ES6.0/libpvr2d.so.1.6.16.3977
/usr/lib/libsrv_init.so.1.6.16.3977
/usr/lib/libOpenVG.so
/usr/lib/ES2.0/libIMGegl.so.1.6.16.3977
/usr/lib/libOpenVG.so.1.6.16.3977
/etc/udev/rules.d/99-bufferclass.rules
/usr/lib/libglslcompiler.so.1.6.16.3977
/usr/lib/libpvr2d.so.1.6.16.3977
/usr/lib/libGLES_CM.so
/usr/bin/ES5.0/pvrsrvinit
/usr/lib/libOpenVGU.so.1.6
/usr/lib/libOpenVG.so.1
/usr/lib/ES5.0/libEGL.so.1.6.16.3977
/usr/lib/ES2.0/libsrv_um.so.1.6.16.3977
/usr/lib/ES6.0/libGLESv2.so.1.6.16.3977
/usr/lib/libOpenVG.so.1.6
/usr/lib/ES5.0/libpvr2d.so.1.6.16.3977
/usr/bin/pvrsrvinit
/usr/lib/ES3.0/libGLES_CM.so.1.6.16.3977
/usr/lib/libGLESv2.so.1
/usr/lib/ES3.0/libOpenVGU.so.1.6.16.3977
